# root/etc/services.d/immich-server/run
#!/bin/sh
exec 2>&1

echo "Starting Immich Server..."
cd /app
exec s6-setuidgid bsd /usr/local/bin/node dist/main.js
